# templates/systemd.service.erb
[Unit]
Description=ConfD Service
After=network.target

[Service]
ExecStart=<%= scope.lookupvar('::confd::installdir') %>/confd
ExecStop=/bin/kill -s TERM $MAINPID

User=<% scope.lookupvar('::confd::user') %>
Group=<% scope.lookupvar('::confd::group') %>

[Install]
WantedBy=multi-user.target
